# DockerConfiguration/supervisor/php.conf
[program:php7-fpm]
command=/usr/sbin/php7-fpm -c /etc/php7/fpm
