# DockerConfiguration/supervisor/postgresql.conf
[program:postgresql]
directory=/usr/local/
command=sh /etc/supervisor/conf.d/postgresql-start.sh
autostart=true
autorestart=true
redirect_stderr=true
